# email-service/supervisord.conf
[supervisord]
nodaemon=true

[program:dropwizard-example]
command=java -jar /app/target/dropwizard-example-0.9.0-SNAPSHOT.jar server /app/example.yml
